# systemd/comfyui-server@.service
[Unit]
Description=ComfyUI image generation server - %i (ROCm)
After=default.target

[Service]
Type=exec
ExecStart=/home/jaegle/local-llm/comfyui-server %i

# ComfyUI and the bridge run inside the distrobox container's cgroup,
# not this unit's, so KillMode= can't reach them.
# Signal the in-container ComfyUI python process and bridge by matching
# their command lines. ExecStopPost SIGKILLs anything left.
ExecStop=-/usr/bin/pkill --signal SIGTERM -f "python main.py.*--port.*8188"
ExecStop=-/usr/bin/pkill --signal SIGTERM -f "comfyui-bridge.py"
ExecStopPost=-/usr/bin/pkill --signal SIGKILL -f "python main.py.*--port.*8188"
ExecStopPost=-/usr/bin/pkill --signal SIGKILL -f "comfyui-bridge.py"
TimeoutStopSec=30

Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
